# userpatches/overlay/etc/systemd/system/pymc-first-boot.service
[Unit]
Description=pyMC Repeater First Boot Setup
After=local-fs.target
Before=sshd.service pymc-repeater.service networking.service
ConditionPathExists=!/etc/pymc-first-boot-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pymc-first-boot.sh
RemainAfterExit=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pymc-first-boot

[Install]
WantedBy=multi-user.target
