# userpatches/overlay/etc/systemd/system/pymc-repeater.service
[Unit]
Description=pyMC LoRa Mesh Repeater Daemon
After=network-online.target pymc-first-boot.service
Wants=network-online.target
ConditionPathExists=/etc/pymc-first-boot-done

[Service]
User=pymc
Group=pymc
Environment=MALLOC_ARENA_MAX=2
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/pymc_repeater/pyMC_Repeater
ExecStartPre=/usr/local/bin/pymc-compile-bytecode.sh
ExecStart=/opt/pymc_repeater/venv/bin/python -O -m repeater.main --config /etc/pymc_repeater/config.yaml
Restart=always
RestartSec=300

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pymc-repeater

ProtectSystem=strict
ReadWritePaths=/etc/pymc_repeater /var/log/pymc_repeater /var/lib/pymc_repeater /opt/pymc_repeater /tmp
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
